# systemd-configs/riing-sleep.timer
[Unit]
Description=Switch to Sleep Config at 3AM

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
