# systemd-configs/riing-sleep.service
[Unit]
Description=Switch Riing Trio to Sleep Config
After=riing-trio-controller.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'ln -sf /etc/riing-trio/sleep.toml /etc/riing-trio/active.toml && systemctl restart riing-trio-controller.service'
