# systemd-configs/riing-work.timer
[Unit]
Description=Switch to Work Config at 10AM

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
